# mysqldump.conf.d/crontab
* * * * * sh /mnt/mysqldump.conf.d/mysqldump.sh >> /mnt/mysqldump.conf.d/log/cron.log
